# scripts/nba-healthcheck.service
[Unit]
Description=NBA production health check

[Service]
Type=oneshot
WorkingDirectory=%h/app
ExecStart=%h/app/scripts/monitor-health.sh
