# scripts/nba-backup.service
[Unit]
Description=NBA Dashboard backup job

[Service]
Type=oneshot
WorkingDirectory=%h/app
Environment=ENVIRONMENT=prod
ExecStart=%h/app/scripts/backup.sh
